# central_hub/node/2-工业项目/nodes.ini
[central_hub_app]
supervisor_username = xyzrobot
supervisor_password = robot2021
supervisor_address = 127.0.0.1:9001
hub_server_address = 127.0.0.1:9004
node_server_address = 127.0.0.1:10001
node_error_url = http://127.0.0.1:7002/api/notify/node_error
dependency_mode = false
boot_sequence = photoneo_no_gui,xvf_no_gui,vision_server
node_monitor_timeout=24

[program:0-auto_start_nodes]
command = python -u /home/xyz/xyz_app/central_hub/script/auto_start_nodes.py
process_name = %(program_name)s
priority = 10
autostart = true
autorestart = false
startsecs = 1
startretries = 1
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = auto_start_nodes
_nodename = auto start all of nodes
_monitor = false
_popup = false
_dependency = 
_description = 自动按顺序启动脚本auto_start_nodes.py中的节点，可以设置延时启动

[program:0-start_all_nodes]
command = python -u /home/xyz/xyz_app/central_hub/script/start_all_nodes.py
process_name = %(program_name)s
priority = 1
autostart = true
autorestart = false
startsecs = 1
startretries = 1
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = start_all_nodes.py
_nodename = start_all_nodes
_monitor = false
_popup = false
_dependency = 
_description = 根据节点启动顺序列表中的节点，使其按顺序启动

[program:1-photoneo]
environment = PHOXI_WITHOUT_DISPLAY=OFF,DISPLAY=":1",PHOXI_CONTROL_PATH=/opt/Photoneo/PhoXiControl-1.4.1
command = /opt/Photoneo/PhoXiControl-1.4.1/bin/PhoXiControl
process_name = %(program_name)s
priority = 10
autostart = false
autorestart = false
startsecs = 5
startretries = 1
stopsignal = TERM
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = photoneo
_nodename = Photo Connector
_monitor = false
_popup = false

[program:1-photoneo_no_gui]
environment = PHOXI_WITHOUT_DISPLAY=ON,PHOXI_CONTROL_PATH=/opt/Photoneo/PhoXiControl-1.4.1
command = /opt/Photoneo/PhoXiControl-1.4.1/bin/PhoXiControl
process_name = %(program_name)s
priority = 10
autostart = false
autorestart = false
startsecs = 5
startretries = 1
stopsignal = TERM
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = photoneo_no_gui
_nodename = Photo Connector
_monitor = false
_popup = false

[program:2-xvf]
environment = PHO_LOG_FILES_DIR=/home/xyz/xyz_log,DISPLAY=":1"
command = /opt/ros/noetic/lib/xyz_vision_flow/xvf_bin
process_name = %(program_name)s
priority = 20
autostart = false
autorestart = false
startsecs = 5
startretries = 1
stopsignal = INT
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = xvf
_nodename = Vision Flow
_monitor = false
_popup = false
_dependency = 
_description = 

[program:2-xvf_no_gui]
environment = PHO_LOG_FILES_DIR=/home/xyz/xyz_log,QT_QPA_PLATFORM="offscreen"
command = /opt/ros/noetic/lib/xyz_vision_flow/xvf_bin --no-gui
process_name = %(program_name)s
priority = 20
autostart = false
autorestart = false
startsecs = 5
startretries = 1
stopsignal = INT
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = xvf_no_gui
_nodename = Vision Flow
_monitor = false
_popup = false

[program:3-vision_server]
command = /usr/bin/xyz_vision_server
process_name = %(program_name)s
priority = 30
autostart = false
autorestart = false
startsecs = 5
startretries = 1
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = vision_server
_nodename = Vision Server
_monitor = false
_popup = false
