# central_hub/node/nodes.ini
[central_hub_app]
supervisor_username = xyzrobot
supervisor_password = robot2020
supervisor_address = 127.0.0.1:9001
hub_server_address = 127.0.0.1:9004
node_server_address = 127.0.0.1:10001
node_error_url = http://127.0.0.1:7002/api/notify/node_error
dependency_mode = false
boot_sequence = robot_config,io_node,robot_node,env_manager,xtf
node_monitor_timeout = 24

[program:0-auto_start_nodes]
command = python -u /home/xyz/xyz_app/central_hub/script/auto_start_nodes.py
process_name = %(program_name)s
priority = 10
autostart = true
autorestart = false
startsecs = 1
startretries = 1
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = auto_start_nodes
_nodename = auto start all of nodes
_monitor = false
_popup = false
_dependency = 
_description = 自动按顺序启动脚本auto_start_nodes.py中的节点，可以设置延时启动

[program:0-data_uploader]
command = /opt/ros/noetic/bin/roslaunch /opt/ros/noetic/share/data_uploader/launch/main.launch upload_cg:=/home/xyz/xyz_app/uploader/uploader_config.yml
process_name = %(program_name)s
priority = 10
autostart = false
autorestart = false
startsecs = 1
startretries = 1
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = data_uploader
_nodename = Data Uploader
_monitor = false
_popup = false
_dependency = 
_description = 

[program:0-hmi_backend]
command = /bin/bash -c "sleep 5 && /usr/bin/xlhb"
priority = 10
autostart = true
autorestart = false
startsecs = 1
startretries = 3
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = hmi_backend
_nodename = HMI Backend Server
_monitor = false
_popup = false
environment = DISPLAY=":1"

[program:1-robot_config]
command = /opt/ros/noetic/bin/roslaunch /home/xyz/xyz_app/app/launch/system_config.launch have_robot:=true
process_name = %(program_name)s
priority = 10
autostart = false
autorestart = false
startsecs = 3
startretries = 2
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = robot_config
_nodename = 机器人配置
_monitor = false
_popup = false

[program:2-io_node]
command = /opt/ros/noetic/bin/rosrun xyz_io io_server -d /home/xyz/xyz_app/app/io/io_server.json
process_name = %(program_name)s
priority = 20
autostart = false
autorestart = false
startsecs = 1
startretries = 2
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = io_node
_nodename = I/O 节点
_monitor = false
_popup = false

[program:3-robot_node]
command = /usr/bin/python -u /home/xyz/xyz_app/central_hub/script/start_robot_node.py
process_name = %(program_name)s
priority = 40
autostart = false
autorestart = false
startsecs = 2
startretries = 1
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = robot_node
_nodename = 机器人节点
_monitor = true
_popup = false

[program:4-env_manager]
command = /opt/ros/noetic/bin/rosrun xyz_env_manager xyz_env_manager_node
process_name = %(program_name)s
priority = 30
autostart = false
autorestart = false
startsecs = 1
startretries = 1
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = env_manager
_nodename = 环境节点
_monitor = false
_popup = false

[program:5-photoneo]
environment = PHOXI_WITHOUT_DISPLAY=OFF,DISPLAY=":1",PHOXI_CONTROL_PATH=/opt/Photoneo/PhoXiControl-1.4.1,PHO_LOG_FILES_DIR=/home/xyz/xyz_log
command = /opt/Photoneo/PhoXiControl-1.4.1/bin/PhoXiControl
process_name = %(program_name)s
priority = 10
autostart = false
autorestart = false
startsecs = 5
startretries = 1
stopsignal = TERM
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = photoneo
_nodename = Photo Connector
_monitor = false
_popup = false

[program:5-photoneo_no_gui]
environment = PHOXI_WITHOUT_DISPLAY=ON,PHOXI_CONTROL_PATH=/opt/Photoneo/PhoXiControl-1.4.1
command = /opt/Photoneo/PhoXiControl-1.4.1/bin/PhoXiControl
process_name = %(program_name)s
priority = 10
autostart = false
autorestart = false
startsecs = 5
startretries = 1
stopsignal = TERM
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = photoneo_no_gui
_nodename = Photo Connector
_monitor = false
_popup = false

[program:6-xvf]
environment = PHO_LOG_FILES_DIR=/home/xyz/xyz_log,DISPLAY=":1"
command = /opt/ros/noetic/lib/xyz_vision_flow/xvf_bin
process_name = %(program_name)s
priority = 50
autostart = false
autorestart = false
startsecs = 5
startretries = 1
stopsignal = INT
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = xvf
_nodename = Vision Flow
_monitor = false
_popup = false
_dependency = 
_description = 

[program:6-xvf_no_gui]
environment = PHO_LOG_FILES_DIR=/home/xyz/xyz_log,QT_QPA_PLATFORM="offscreen"
command = /opt/ros/noetic/lib/xyz_vision_flow/xvf_bin --no-gui
process_name = %(program_name)s
priority = 50
autostart = false
autorestart = false
startsecs = 5
startretries = 1
stopsignal = INT
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = xvf_no_gui
_nodename = Vision Flow
_monitor = false
_popup = false
_dependency = 
_description = 

[program:6-xtf]
command = /opt/ros/noetic/bin/rosrun xtf_ros xtf_node
process_name = %(program_name)s
priority = 60
autostart = false
autorestart = false
startsecs = 2
startretries = 1
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = xtf
_nodename = 执行任务节点
_monitor = false
_popup = false
environment = PATH="/usr/local/bin:/usr/local/cuda/bin:/opt/ros/noetic/bin:/home/xyz/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin:/home/xyz/xyz_app/software/build/bin",
	LD_LIBRARY_PATH="/usr/local/lib:/opt/vcpkg/installed/x64-linux/lib:/home/xyz/xyz-studio-max/1.9.0-pre.17/lib:/opt/ros/noetic/lib:/usr/local/lib:/usr/local/cuda/lib:/opt/MVS/lib/64:/opt/MVS/lib/32:/opt/MVS/lib/64:/opt/MVS/lib/32:/opt/MVS/lib/64:/opt/MVS/lib/32",
	PYTHONPATH="/home/xyz/xyz-studio-max/1.9.0-pre.17/lib/python/dist-packages:/opt/ros/noetic/lib/python3/dist-packages",
	XYZ_SERVER_ID="2204",
	XYZ_SERVER_NAME="dapeng_station_0",
	PHOXI_WITHOUT_DISPLAY=ON,
	PHOXI_CONTROL_PATH=/opt/Photoneo/PhoXiControl-1.4.1,
	PHO_LOG_FILES_DIR=/home/xyz/xyz_log

[program:7-mix_palletize]
command = /opt/ros/noetic/bin/roslaunch /home/xyz/xyz_app/app/launch/mix_palletize.launch
process_name = %(program_name)s
priority = 70
autostart = false
autorestart = false
startsecs = 2
startretries = 1
user = xyz
redirect_stderr = true
stdout_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
stdout_capture_maxbytes = 10MB
stdout_events_enabled = false
stderr_logfile = /home/xyz/xyz_log/central_hub/nodes/%(program_name)s.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 2
stderr_capture_maxbytes = 10MB
stderr_events_enabled = false
_nodeid = mix_palletize
_nodename = Mix Palletize
_monitor = false
_popup = false

